# system_files/lib/systemd/system/btmidi.service
[Unit]
Description=MIDI Bluetooth connect
After=bluetooth.target
Requires=bluetooth.target

[Service]
ExecStartPre=/bin/sleep 15
Type=simple
User=root
Group=root
WorkingDirectory=/home/pi
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=btmidi
Restart=always
ExecStart=/usr/bin/btmidi-server -n "RPi_LED_BTMIDI-Server"

[Install]
WantedBy=multi-user.target
